# scripts/katok-watcher.service
[Unit]
Description=KakaoTalk File Watcher
After=network.target

[Service]
Type=simple
ExecStart=/home/ppak/clawd/scripts/katok-watcher-daemon.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
